# docker/root/etc/services.d/api/run
#!/usr/bin/with-contenv sh

if [ "$NODE_ENV" = "development" ]; then
    command="watch";
else
    command="start";
fi

cd /app

eval npm run $command;
